# debian/postinst
#!/bin/sh

update-rc.d flussonic defaults 20 80 >/dev/null



